# Sync/Linux/systemd/adhs-lernpfad-sync.service
[Unit]
Description=ADHS-Lernpfad mit Obsidian synchronisieren
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=ADHS_SYNC_NONINTERACTIVE=1
ExecStart=%h/.local/bin/adhs-lernpfad-sync
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
