# Sync/Linux/systemd/adhs-lernpfad-sync.service
[Unit]
Description=ADHS-Lernpfad von GitHub nach Obsidian synchronisieren
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/adhs-lernpfad-sync.env
ExecStart=%h/.local/bin/adhs-lernpfad-sync
Nice=10
